# Library/LaunchAgents/io.github.henrylinyy.qwable.warmup.quality.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.henrylinyy.qwable.warmup.quality</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yourname/Documents/qwable-agent-gateway-m5/scripts/warmup_fusion_quality.sh</string>
        <string>--quiet</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yourname/Documents/qwable-agent-gateway-m5</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yourname/.lmstudio/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OLLAMA_BASE_URL</key>
        <string>http://127.0.0.1:1234/v1</string>
    </dict>

    <!-- Run weekdays (Mon=1 ... Fri=5) at 08:00 local time -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/yourname/Library/Logs/qwable-warmup-quality.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yourname/Library/Logs/qwable-warmup-quality.err.log</string>
</dict>
</plist>
